# min-stack/bundle/s6-rc.d/sink/run
#!/command/execlineb -P
with-contenv
s6-envdir /work/env/_shared
s6-envdir /work/env/sink
fdmove -c 2 1
/work/bin/sink
