# min-stack/bundle/s6-rc.d/db/run
#!/command/execlineb -P
with-contenv
s6-envdir /work/env/db
fdmove -c 2 1
/work/bin/db
