# min-stack/bundle/s6-rc.d/ender/run
#!/command/execlineb -P
with-contenv
s6-envdir /work/env/ender
fdmove -c 2 1
/work/bin/ender
