# min-stack/bundle/s6-rc.d/storage/run
#!/command/execlineb -P
with-contenv
s6-envdir /work/env/storage
fdmove -c 2 1
/work/bin/storage
